# to_install/display-cameras-restart.service
[Unit]
Description=Display cameras daily restart

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart display-camera@ch*

